# contrib/init/doriancoind.conf
description "Doriancoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env DORIANCOIND_BIN="/usr/bin/doriancoind"
env DORIANCOIND_USER="doriancoin"
env DORIANCOIND_GROUP="doriancoin"
env DORIANCOIND_PIDDIR="/var/run/doriancoind"
# upstart can't handle variables constructed with other variables
env DORIANCOIND_PIDFILE="/var/run/doriancoind/doriancoind.pid"
env DORIANCOIND_CONFIGFILE="/etc/doriancoin/doriancoin.conf"
env DORIANCOIND_DATADIR="/var/lib/doriancoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # doriancoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$DORIANCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run doriancoind."
        echo "The setting must appear in $DORIANCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Doriancoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$DORIANCOIND_PIDDIR"
    chmod 0755 "$DORIANCOIND_PIDDIR"
    chown $DORIANCOIND_USER:$DORIANCOIND_GROUP "$DORIANCOIND_PIDDIR"
    chown $DORIANCOIND_USER:$DORIANCOIND_GROUP "$DORIANCOIND_CONFIGFILE"
    chmod 0660 "$DORIANCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$DORIANCOIND_PIDFILE" \
    --chuid $DORIANCOIND_USER:$DORIANCOIND_GROUP \
    --exec "$DORIANCOIND_BIN" \
    -- \
    -pid="$DORIANCOIND_PIDFILE" \
    -conf="$DORIANCOIND_CONFIGFILE" \
    -datadir="$DORIANCOIND_DATADIR" \
    -disablewallet \
    -daemon

